# user-services/dotfiles-deploy.path
[Path]
PathChanged=%h/.dotfiles

[Install]
WantedBy=default.target
