# user-services/dotfiles-deploy.service
[Unit]
Description=Watch and deploy dotfiles using bombadil

[Service]
ExecStart=bombadil watch -p {{profile}}

[Install]
WantedBy=default.target
